# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── PIPELINE A — FAST ─────────────────────────────────────────────────────
[program:fast_pipeline]
command=python /app/fast_pipeline.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/fast_pipeline.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/fast_pipeline_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; env vars inherited from container (docker-compose.yml) — do NOT hardcode here

; ── PIPELINE B — RETRY ────────────────────────────────────────────────────
[program:retry_pipeline]
command=python /app/retry_pipeline.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/retry_pipeline.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/retry_pipeline_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; env vars inherited from container (docker-compose.yml) — do NOT hardcode here
