# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── PIPELINE A — FAST ─────────────────────────────────────────────────────
[program:fast_pipeline]
command=python /app/fast_pipeline.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/fast_pipeline.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/fast_pipeline_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; BUG-16 FIX: all key=value pairs on ONE line — multi-line with leading spaces is invalid
environment=FAST_NUM_WORKERS="15",FAST_BATCH_SIZE="25",FAST_TIMEOUT_MS="3000",FAST_STAGGER_S="0.5"

; ── PIPELINE B — RETRY ────────────────────────────────────────────────────
[program:retry_pipeline]
command=python /app/retry_pipeline.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/retry_pipeline.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/retry_pipeline_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; BUG-16 FIX: single-line environment
environment=RETRY_NUM_WORKERS="5",RETRY_BATCH_SIZE="10",RETRY_TIMEOUT_MS="15000",RETRY_STAGGER_S="1.0"
